# crontab
0 3 * * * root certbot renew >> /var/log/letsencrypt.log 2>&1
